# deploy/iperf-speedtest.service
[Unit]
Description=iperf3 speedtest collector (Init7)
After=network-online.target
Wants=network-online.target

[Service]
# Path is rewritten by deploy/install.sh to wherever the repo is checked out.
ExecStart=/usr/bin/python3 /opt/iperf-speedtest/collector.py
Environment=HTTP_PORT=8089
Environment=DATA_DIR=/var/lib/iperf-speedtest
Environment=IPERF_HOST=speedtest.init7.net
Environment=IPERF_PORT=5202
Environment=IPERF_DURATION=10
Environment=IPERF_OMIT=2
# Runs once a day at a random time in this local window (kind to the test server).
Environment=TZ=Europe/Zurich
Environment=DAILY_WINDOW_START_HOUR=2
Environment=DAILY_WINDOW_END_HOUR=5
# IPERF_PARALLEL is left unset → defaults to 16 (see collector.py).
DynamicUser=yes
StateDirectory=iperf-speedtest
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
